# com.mikehu.cmdr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.mikehu.cmdr</string>
    <key>ProgramArguments</key>
    <array>
      <string>__CMDR_BIN__</string>
      <string>start</string>
      <string>-f</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/cmdr.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cmdr.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
      <key>CMDR_CODE_DIR</key>
      <string>~/Code</string>
      <key>CMDR_OLLAMA_URL</key>
      <string>https://ollama.106source.ca</string>
      <key>CMDR_OLLAMA_MODEL</key>
      <string>gemma4:e4b</string>
    </dict>
  </dict>
</plist>
